# docker/etc/services.d/worker-email/run
#!/command/execlineb -P

/var/www/stargaze/bin/console messenger:consume async -q --memory-limit=128M --queues email_notifications
